# flock-sensor/flock-sensor.service
[Unit]
Description=Flock Venue Sensor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/flock-sensor
ExecStart=/usr/bin/python3 /home/pi/flock-sensor/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
